# packages/luci-app-uestc-authclient/root/etc/uci-defaults/99-uestc-authclient-ruijie-migrate
#!/bin/sh

changed=0

if ! uci -q get uestc_authclient.ct_ruijie >/dev/null 2>&1; then
    old_user="$(uci -q get uestc_authclient.ct.auth_username || true)"
    old_pass="$(uci -q get uestc_authclient.ct.auth_password || true)"
    old_iface="$(uci -q get uestc_authclient.ct.listen_interface || true)"
    old_lm="$(uci -q get uestc_authclient.ct.lm_enabled || true)"

    [ -n "$old_iface" ] || old_iface="wan"
    [ -n "$old_lm" ] || old_lm="1"

    uci set uestc_authclient.ct_ruijie='session'
    uci set uestc_authclient.ct_ruijie.enabled='0'
    uci set uestc_authclient.ct_ruijie.lm_enabled="$old_lm"
    uci set uestc_authclient.ct_ruijie.auth_type='ct_ruijie'
    uci set uestc_authclient.ct_ruijie.auth_host='110.184.24.61'
    uci set uestc_authclient.ct_ruijie.listen_interface="$old_iface"
    uci set uestc_authclient.ct_ruijie.listen_check_interval='30'
    uci -q delete uestc_authclient.ct_ruijie.listen_hosts
    uci add_list uestc_authclient.ct_ruijie.listen_hosts='223.5.5.5'
    uci add_list uestc_authclient.ct_ruijie.listen_hosts='119.29.29.29'
    uci set uestc_authclient.ct_ruijie.schedule_enabled='0'
    uci set uestc_authclient.ct_ruijie.schedule_start='3'
    uci set uestc_authclient.ct_ruijie.schedule_end='4'
    uci set uestc_authclient.ct_ruijie.log_rdays='7'

    [ -n "$old_user" ] && uci set uestc_authclient.ct_ruijie.auth_username="$old_user"
    [ -n "$old_pass" ] && uci set uestc_authclient.ct_ruijie.auth_password="$old_pass"
    changed=1
fi

if [ "$(uci -q get uestc_authclient.ct.auth_type || true)" = "ct" ] &&
   [ "$(uci -q get uestc_authclient.ct.auth_host || true)" = "110.184.24.61" ]; then
    uci set uestc_authclient.ct.auth_host='172.25.249.64'
    changed=1
fi

[ "$changed" -eq 1 ] && uci commit uestc_authclient
exit 0
