# packages/luci-app-uestc-authclient/root/etc/uci-defaults/99-uestc-authclient-ruijie-migrate
#!/bin/sh

changed=0
telecom_entry_host="172.25.249.64"

normalize_telecom_section() {
    section="$1"

    uci -q get "uestc_authclient.$section" >/dev/null 2>&1 || return

    auth_type="$(uci -q get "uestc_authclient.$section.auth_type" || true)"
    case "$auth_type" in
        ""|ct|ct_ruijie|qsh-telecom-ruijie)
            if [ "$auth_type" != "ct" ]; then
                uci set "uestc_authclient.$section.auth_type=ct"
                changed=1
            fi
            ;;
    esac

    auth_host="$(uci -q get "uestc_authclient.$section.auth_host" || true)"
    case "$auth_host" in
        ""|110.184.24.61)
            uci set "uestc_authclient.$section.auth_host=$telecom_entry_host"
            changed=1
            ;;
    esac
}

if uci -q get "uestc_authclient.ct_ruijie" >/dev/null 2>&1 &&
   ! uci -q get "uestc_authclient.ct" >/dev/null 2>&1; then
    if uci rename "uestc_authclient.ct_ruijie=ct" 2>/dev/null; then
        changed=1
    fi
fi

if uci -q get "uestc_authclient.qsh_telecom_ruijie" >/dev/null 2>&1 &&
   ! uci -q get "uestc_authclient.ct" >/dev/null 2>&1; then
    if uci rename "uestc_authclient.qsh_telecom_ruijie=ct" 2>/dev/null; then
        changed=1
    fi
fi

normalize_telecom_section "ct"
normalize_telecom_section "ct_ruijie"
normalize_telecom_section "qsh_telecom_ruijie"

[ "$changed" -eq 1 ] && uci commit uestc_authclient
exit 0
